# s3_backup_tool/rootfs/etc/services.d/s3-backup/run
#!/usr/bin/with-contenv bashio

# Starte das Add-on Hauptskript
exec /run.sh


